# arches/install/supervisor_celery_setup/my_proj_name-celerybeat.conf
; ================================
;  celery beat supervisor
; ================================

[program:celerybeat]
	command=/path/to/virtualenv/bin/celery beat -A [app].celery --schedule=/tmp/celerybeat-schedule --loglevel=INFO --pidfile=/tmp/celerybeat.pid
	directory=/absolute/path/to/my_proj

    user=%(ENV_USER)s
    numprocs=1
    stdout_logfile=/var/log/celery/beat.log
    stderr_logfile=/var/log/celery/beat.log
    autostart=true
    autorestart=true
    startsecs=10

    ; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
    stopasgroup=true

    ; if rabbitmq is supervised, set its priority higher
    ; so it starts first
    priority=999